# ansible/files/speak-to-me/ddclient/containers_systemd/ddclient.container
[Unit]
Description=ddclient container

[Container]
Image=docker.io/linuxserver/ddclient:3.11.2
Volume=/home/ddclient/config:/config:Z
Environment="TZ=Europe/Warsaw"
Secret=DDCLIENT_PASSWORD,type=env

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
